# packages/root-compatibility-enviroment/etc/profile.d/n64-root.sh
# Declare the /etc/n64 environment
export ROOT=/etc/n64
